# Library/LaunchAgents/kardia.notes.auto_commit.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>kardia.notes.auto-commit.plist</string>

    <key>ProcessType</key>
    <string>Interactive</string>

    <key>ProgramArguments</key>
    <array>
      <string>/Users/kardia/notes/auto-commit.sh</string>
    </array>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>StandardErrorPath</key>
<string>/tmp/launchdebug.err</string>
<key>StandardOutPath</key>
<string>/tmp/launchdebug.out</string>

  </dict>
</plist>
